# src/99-mw_initdb.conf
[group:mediawiki]
programs=mw

[program:mw]
command = echo 0
process_name=%(program_name)s
priority=4
startsecs=0
autostart = true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


